# ansible/templates/fix-box.service.j2
[Unit]
Before=network-pre.target
Wants=network-pre.target

[Service]
ExecStart=/usr/local/bin/fix-box.sh

[Install]
WantedBy=default.target
